# commonConf/v0.0.2-V0.0.2-merge-all_linux_system/nginx/kylin-x86/service/nginx.service
[Unit]
Description=nginx
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/nginx -c /usr/local/nginx/nginx.conf	
PrivateTmp=true

[Install]
WantedBy=multi-user.target
